# node-exporter-sytemd/files/node-exporter.service
[Unit]
Description=Node Exporter

[Service]
ExecStart=/opt/node-exporter/node-exporter

[Install]
WantedBy=multi-user.target
